# zabbix-postfix.cron
# Zabbix check for postfix-mailcow
*/5 * * * * /etc/zabbix/scripts/zabbix_postfix.sh 1>/dev/null 2>/dev/null
