# zabbix-postfix.timer
[Unit]
Description=update zabbix-postfix-mailcow every 30 seconds

[Timer]
OnActiveSec=30
Persistent=true

[Install]
WantedBy=timers.target

